# packaging/systemd/spanel-api.service
[Unit]
Description=SimpleHost Panel API
After=network-online.target postgresql@shp.service
Wants=network-online.target
Requires=postgresql@shp.service

[Service]
Type=simple
User=spanel
Group=spanel
WorkingDirectory=/opt/simplehost/spanel/current
EnvironmentFile=/etc/spanel/api.env
ExecStart=/usr/bin/node /opt/simplehost/spanel/current/apps/api/dist/index.js
Restart=always
RestartSec=5
StateDirectory=spanel
LogsDirectory=spanel

[Install]
WantedBy=multi-user.target
