# packaging/systemd/spanel-web.service
[Unit]
Description=SimpleHost Panel Web
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=spanel
Group=spanel
WorkingDirectory=/opt/simplehost/repos/simplehost-panel
EnvironmentFile=/etc/spanel/web.env
ExecStart=/usr/bin/node /opt/simplehost/repos/simplehost-panel/apps/web/dist/index.js
Restart=always
RestartSec=5
StateDirectory=spanel
LogsDirectory=spanel

[Install]
WantedBy=multi-user.target
